# deploy/packaging/prerm
#!/bin/bash
set -e
systemctl stop qf-agent || true
systemctl disable qf-agent || true
